# summary-service/ops/ycare-summary-api.service
[Unit]
Description=YCare HMS — Summary Service (HTTP API)
Documentation=https://internal.example.com/docs/summary-service
# Run after Postgres and Redis are up; the API will not start successfully
# if either is unreachable, so we order strictly.
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=simple
User=ycare-summary
Group=ycare-summary

# The compiled app is at /opt/ycare-summary/dist/index.js.
# The --mode=api flag routes to the Express HTTP server.
ExecStart=/usr/bin/node /opt/ycare-summary/dist/index.js --mode=api
ExecReload=/bin/kill -HUP $MAINPID

# Environment
EnvironmentFile=/etc/ycare-summary/env

# Resource limits (API is request/response; can spike under load)
LimitNOFILE=65536
MemoryMax=512M
CPUQuota=100%

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=10

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native

# Read-only paths
ReadOnlyPaths=/opt/ycare-summary
ReadOnlyPaths=/etc/ycare-summary

# Writable paths
RuntimeDirectory=ycare-summary
RuntimeDirectoryMode=0750
StateDirectory=ycare-summary
StateDirectoryMode=0750

# Logging
StandardOutput=append:/var/log/ycare-summary/api.log
StandardError=append:/var/log/ycare-summary/api.error.log

[Install]
WantedBy=multi-user.target
